# backup_selfhosted/backup_selfhosted.service
[Unit]
Description=Backup Data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup_selfhosted.sh

[Install]
WantedBy=multi-user.target